# hazel-chat-webhook.service
[Unit]
Description=Hazel Dashboard Chat Webhook
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/hazel-chat-webhook
ExecStart=/home/openclaw/hazel-chat-webhook/venv/bin/python server.py
Restart=always
RestartSec=5
EnvironmentFile=/home/openclaw/hazel-chat-webhook/.env

[Install]
WantedBy=multi-user.target
